# install/systemd/event/event_is.service
[Unit]
Description=Event service for Intersvyaz intercom syslog messages
Documentation=https://github.com/rosteleset/SmartYard-Server/tree/main/install
After=network.target

[Service]
SyslogIdentifier=event_is
Type=simple
WorkingDirectory=/opt/rbt/server/services/event/
ExecStart=/usr/bin/node /opt/rbt/server/services/event/index.js is
RestartSec=10
Restart=always
User=rbt
Group=rbt
LimitCORE=infinity
KillMode=process
StandardOutput=null
#StandardOutput=append:/var/log/event_service/event_is.log
StandardError=append:/var/log/event_service/event_is.error.log

[Install]
WantedBy=multi-user.target
